# nginx_docker/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=10
process_name=%(program_name)s
redirect_stderr=true
stderr_logfile=/var/log/nginx_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backup=5
stdout_logfile=/var/log/nginx_access.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backup=5

